# com.zeev.learningos.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  תזמון אוטומטי (launchd) — ריצה שבועית מעמיקה, ראשון ב-08:00.
  התקנה:
    1. ערוך את הנתיבים למטה אם הזזת את התיקייה או node.
    2. cp com.zeev.learningos.plist ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/com.zeev.learningos.plist
  להסרה:  launchctl unload ~/Library/LaunchAgents/com.zeev.learningos.plist
  (המחשב חייב להיות דלוק בזמן הריצה; אם ישן — launchd יריץ עם ההתעוררות.)
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zeev.learningos</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/sarit/Desktop/learning-os/scripts/auto-refresh.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/sarit/Desktop/learning-os</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/sarit/Desktop/learning-os/digest.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/sarit/Desktop/learning-os/digest.log</string>
</dict>
</plist>
